# rootfs/etc/services.d/ozric/run
#!/usr/bin/with-contenv bashio

exec /ozric/Ozric.Dashboard
